# gitbucket_home/service/gitbucket_daemonized.service
[Unit]
Description=webserver Daemon

[Service]
ExecStart=/usr/bin/java -jar /usr/lib/gitbucket.war
User=root

[Install]
WantedBy=multi-user.target
